# deploy/vortex.service
# ══════════════════════════════════════════════════════════════════════════════
# VORTEX Chat — systemd Service Unit
# ══════════════════════════════════════════════════════════════════════════════
#
# Installation:
#   1. Copy to /etc/systemd/system/vortex.service
#   2. Edit paths, user, and environment as needed
#   3. systemctl daemon-reload
#   4. systemctl enable --now vortex
#
# Management:
#   systemctl status vortex      # Check status
#   systemctl restart vortex     # Restart
#   journalctl -u vortex -f      # Follow logs
#
# ══════════════════════════════════════════════════════════════════════════════

[Unit]
Description=Vortex Chat — Decentralized P2P Messenger
Documentation=https://github.com/your-org/vortex
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=vortex
Group=vortex
WorkingDirectory=/opt/vortex

# ── Environment ──────────────────────────────────────────────────────────────
Environment=ENVIRONMENT=production
Environment=HOST=0.0.0.0
Environment=PORT=9000
Environment=LOG_FORMAT=json
Environment=LOG_LEVEL=INFO
EnvironmentFile=-/opt/vortex/.env

# ── Execution ────────────────────────────────────────────────────────────────
ExecStart=/opt/vortex/.venv/bin/python -m gunicorn \
    app.main:app \
    --worker-class uvicorn.workers.UvicornWorker \
    --workers 1 \
    --bind 0.0.0.0:9000 \
    --timeout 120 \
    --graceful-timeout 30 \
    --access-logfile - \
    --error-logfile -

ExecReload=/bin/kill -HUP $MAINPID

# ── Graceful shutdown ────────────────────────────────────────────────────────
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# ── Restart policy ───────────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# ── Security hardening ───────────────────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes

# Allow write to app directory
ReadWritePaths=/opt/vortex/uploads
ReadWritePaths=/opt/vortex/keys
ReadWritePaths=/opt/vortex/certs
ReadWritePaths=/opt/vortex/logs
ReadWritePaths=/opt/vortex/vortex.db
ReadWritePaths=/opt/vortex/vortex.db-wal
ReadWritePaths=/opt/vortex/vortex.db-shm
ReadWritePaths=/opt/vortex/.env

# ── Resource limits ──────────────────────────────────────────────────────────
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
CPUQuota=200%

# ── Logging ──────────────────────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vortex

[Install]
WantedBy=multi-user.target
